# deploy/ai-radio-dynamic-daemon.service
[Unit]
Description=AI Radio Dynamic Generation Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=martin
Group=martin
WorkingDirectory=/Users/martin/Documents/GitHub/AI-Radio
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/Users/martin/Documents/GitHub/AI-Radio/.env
EnvironmentFile=-/Users/martin/Documents/GitHub/AI-Radio/.env.dynamic-daemon
ExecStart=/Users/martin/Documents/GitHub/AI-Radio/.venv/bin/python /Users/martin/Documents/GitHub/AI-Radio/dynamic_daemon.py --host 127.0.0.1 --port 8787
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
